# build-astroarch/systemd/astroarch-onboarding.timer
[Unit]
Description=Timer astroarch-onboarding

[Timer]
OnBootSec=1min
Unit=astroarch-onboarding.service

[Install]
WantedBy=timers.target
