# build-astroarch/clear-install-astroarch.timer
[Unit]
Description=Timer clear install astroarch

[Timer]
OnBootSec=10min
Unit=clear-install-astroarch.service

[Install]
WantedBy=timers.target
